# netbird-client/rootfs/etc/s6-overlay/s6-rc.d/netbird-client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Netbird client
# ==============================================================================

declare -a options

options+=(-F)

# Setup keys, if defined
SETUP_KEY="$(bashio::config 'SETUP_KEY')"
if [ -z "$SETUP_KEY" ]; then
    bashio::log.warning "No setup key defined. Remember to disable Peer login expiration!"
else
    bashio::log.info "Setup key found, using $SETUP_KEY"
    options+=(--setup-key "${SETUP_KEY}")
fi

# Config file
CONFIG_FILE="$(bashio::config 'CONFIG_FILE')"
options+=(--config "${CONFIG_FILE}")

ADMIN_URL="$(bashio::config 'ADMIN_URL')"
options+=(--admin-url "${ADMIN_URL}")

MANAGEMENT_URL="$(bashio::config 'MANAGEMENT_URL')"
options+=(--management-url "${MANAGEMENT_URL}")

LOG_LEVEL="$(bashio::config 'LOG_LEVEL')"
options+=(--log-level "${LOG_LEVEL}")

bashio::log.info "Netbird options:"
bashio::log.info "${options[@]}"

# Run netbird
bashio::log.info "Starting Netbird client..."
exec netbird up "${options[@]}"
