# netbird-client/rootfs/etc/s6-overlay/s6-rc.d/netbird-client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Netbird client
# ==============================================================================

# Declare variables
declare -a options

options+=(-F)

SETUP_KEY="$(bashio::config 'SETUP_KEY')"
if [ -z "$SETUP_KEY" ]; then
    bashio::log.warning "No setup key defined. Remember to disable Peer login expiration!"
else
    bashio::log.info "Setup key found, using $SETUP_KEY"
    options+=(--setup-key "${SETUP_KEY}")
fi

CONFIG_FILE="$(bashio::config 'SETUP_KEY')"
bashio::log.info "Starting Netbird client..."

## Run netbird
exec netbird up -F "${options[@]}"
