# netbird-client/rootfs/etc/s6-overlay/s6-rc.d/netbird-client/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the NetBird client
# ==============================================================================

declare -a options

options+=(-F)

# Config file
if bashio::config.has_value "config_file"; then
    CONFIG_FILE="$(bashio::config "config_file")"
else
    CONFIG_FILE="/config/netbird/config.json"
fi
options+=(--config "${CONFIG_FILE}")

if [ -e "${CONFIG_FILE}" ]; then
    bashio::log.info "Config file found at ${CONFIG_FILE}. Ignoring other options"
else
    if bashio::config.has_value "setup_key"; then
        options+=(--setup-key "$(bashio::config 'setup_key')")
    else
        bashio::log.warning "No setup key defined. Falling back to login-url (See later in this log)"
        bashio::log.info "Remember to disable peer login expiration in the admin panel"
    fi

    if bashio::config.has_value "admin_url"; then
        options+=(--admin-url "$(bashio::config 'admin_url')")
    fi

    if bashio::config.has_value "management_url"; then
        options+=(--management-url "$(bashio::config 'management_url')")
    fi

    LOG_LEVEL="$(bashio::config 'log_level')"
    options+=(--log-level "${LOG_LEVEL}")
fi

bashio::log.info "Netbird version: $(netbird version)"
bashio::log.info "NetBird options: ${options[@]}"

# Run netbird
bashio::log.info "Starting NetBird client..."
exec netbird up "${options[@]}"
