# netbird-client/rootfs/etc/s6-overlay/s6-rc.d/netbird-client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Netbird client
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare -a options

SETUP_KEY="$(bashio::config 'SETUP_KEY')"
if [ -z "$SETUP_KEY" ]; then
    bashio::warning "No setup key defined. Remember to disable Peer login expiration!"
else
    bashio::info "Setup key found, using $SETUP_KEY"
    options+=(--setup-key "${SETUP_KEY}")
fi

CONFIG_FILE="$(bashio::config 'SETUP_KEY')"

bashio::log.info "Hello..."

## Run your program
#exec /usr/bin/my_program
tail -f /dev/null